# foam-n8n-implementation/docker/supervisord.conf
# =============================================================================
# FOAM N8N Multi-LLM Orchestration - Supervisord Configuration
# =============================================================================
# Manages all services within the unified container:
# - PostgreSQL
# - Redis
# - Ollama
# - N8N
# =============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
logfile_maxbytes=50MB
logfile_backups=3
loglevel=info

# =============================================================================
# PostgreSQL - Database Server
# =============================================================================

[program:postgresql]
command=/usr/lib/postgresql/14/bin/postgres -D /var/lib/postgresql/data
user=postgres
autostart=true
autorestart=true
priority=10
startsecs=10
startretries=3
stopwaitsecs=60
stdout_logfile=/var/log/foam/postgresql.log
stderr_logfile=/var/log/foam/postgresql_error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3

# =============================================================================
# Redis - Queue and Cache Server
# =============================================================================

[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf
user=redis
autostart=true
autorestart=true
priority=20
startsecs=5
startretries=3
stopwaitsecs=30
stdout_logfile=/var/log/foam/redis.log
stderr_logfile=/var/log/foam/redis_error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3

# =============================================================================
# Ollama - Local LLM Inference Server
# =============================================================================

[program:ollama]
command=/usr/local/bin/ollama serve
user=root
autostart=true
autorestart=true
priority=30
startsecs=10
startretries=3
stopwaitsecs=60
environment=OLLAMA_HOST="0.0.0.0",OLLAMA_PORT="11434"
stdout_logfile=/var/log/foam/ollama.log
stderr_logfile=/var/log/foam/ollama_error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile_backups=3

# =============================================================================
# N8N - Workflow Automation Engine
# =============================================================================

[program:n8n]
command=/usr/bin/n8n start
user=root
autostart=true
autorestart=true
priority=40
startsecs=30
startretries=5
stopwaitsecs=60
environment=HOME="/root",N8N_USER_FOLDER="/home/node/.n8n"
stdout_logfile=/var/log/foam/n8n.log
stderr_logfile=/var/log/foam/n8n_error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5

# =============================================================================
# Service Groups
# =============================================================================

[group:database]
programs=postgresql,redis
priority=10

[group:ai]
programs=ollama
priority=20

[group:workflow]
programs=n8n
priority=30

# =============================================================================
# Event Listeners
# =============================================================================

[eventlistener:process_exit]
command=/bin/bash -c 'while read line; do echo "Process event: $line" >> /var/log/foam/supervisor_events.log; done'
events=PROCESS_STATE_EXITED,PROCESS_STATE_FATAL
buffer_size=100

# =============================================================================
# Unix HTTP Server (for supervisorctl)
# =============================================================================

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
